# wsgate/conf/wsgate.service.in
[Unit]
Description=FreeRDP WebAccess Gateway
After=syslog.target network.target

[Service]
Type=forking
User=wsgate
Group=wsgate
LimitCore=infinity
PIDFile=/var/run/wsgate/wsgate.pid
Environment=WSGATE_CONFIG=/etc/wsgate.ini
EnvironmentFile=-/etc/sysconfig/wsgate
ExecStart=/usr/sbin/wsgate -c $WSGATE_CONFIG $WSGATE_OPTIONS
ExecStartPre=_pkglibexecdir_/keygen.sh
ExecStartPre=/bin/mkdir -p /var/run/wsgate
ExecStartPre=/bin/chown wsgate.wsgate /var/run/wsgate
ExecReload=/bin/kill -HUP $MAINPID
PermissionsStartOnly=true

[Install]
WantedBy=multi-user.target
